# deploy/cloudflared.service
[Unit]
Description=Cloudflare Tunnel
After=network.target

[Service]
Type=simple
User=coke
ExecStart=/usr/local/bin/cloudflared tunnel --no-autoupdate --url http://localhost:5678
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
