# deploy/n8n.service
[Unit]
Description=n8n workflow automation
After=network.target

[Service]
Type=simple
User=coke
WorkingDirectory=/home/coke/bot-recetas
EnvironmentFile=/home/coke/bot-recetas/.env
Environment=NODE_OPTIONS=--dns-result-order=ipv4first
Environment=N8N_BLOCK_ENV_ACCESS_IN_NODE=false
Environment=GENERIC_TIMEZONE=Europe/Madrid
ExecStart=/home/coke/.nvm/versions/node/v22.23.0/bin/n8n start
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
